# pro-v620/fan-control/gpu-fan-control@.service
[Unit]
Description=GPU fan control (%i) — Radeon Pro V620 cooler, GPU-temp driven
Documentation=https://github.com/Marchah/Proxmox/tree/main/pro-v620/fan-control
# Needs the nct6687 (writable PWM) module and amdgpu bound to the card.
# One instance per GPU/cooler; the per-instance env pins the GPU (by PCI address)
# and its fan pwm channel: gpu-fan-control@blower, gpu-fan-control@arctic.
Wants=systemd-modules-load.service
After=systemd-modules-load.service
ConditionPathExistsGlob=/sys/class/hwmon/hwmon*/name

[Service]
Type=simple
EnvironmentFile=-/etc/gpu-fan-control-%i.env
ExecStart=/usr/local/sbin/gpu-fan-control
# Failsafe (chosen): on ANY stop/crash hand this instance's fan back to the
# BIOS/SIO auto curve. ExecStopPost also covers SIGKILL, where the script's EXIT
# trap cannot run. Reads FAN_PWM_CHANNEL from the per-instance env.
ExecStopPost=/usr/local/sbin/gpu-fan-control --failsafe
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
